# initrd/rootfs/var/service/dhcp/run
#!/bin/sh
exec > >(logger -t dhcp) 2>&1
exec udhcpc -f -i eth0 -s /etc/udhcp/simple.script
