# ramfs/rootfs/etc/service/cron/run
#!/bin/sh
exec 2>&1 > >(logger -t cron)
exec crond -fS
