# initrd/rootfs/etc/inittab
# See https://github.com/brgl/busybox/blob/master/examples/inittab

# Init
::sysinit:/etc/runit/1

# Services
::respawn:/etc/runit/2

# Console
::askfirst:-/bin/sh
# Replace with this if you only want to allow SSH login
#::askfirst:echo "Login over SSH"

# Shutdown
::shutdown:/etc/runit/3
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
