# initrd/rootfs/etc/inittab
# See https://github.com/brgl/busybox/blob/master/examples/inittab

# Init
::sysinit:/etc/runit/1

# Services
::respawn:/etc/runit/2

# Console
# Use fakelogin if you only want to allow SSH login
#::askfirst:/bin/fakelogin
::askfirst:-/bin/sh

# Shutdown
::shutdown:/etc/runit/3
::restart:/bin/init
::ctrlaltdel:/bin/reboot
