# initrd/rootfs/etc/inittab
# BusyBox init configuration
# See https://github.com/brgl/busybox/blob/master/examples/inittab

# Init
::sysinit:/etc/runit/1

# Services
::respawn:/etc/runit/2

# Shutdown
::shutdown:/etc/runit/3
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot

# Console
::askfirst:-/bin/sh
