# initrd/rootfs/etc/service/ssh/run
#!/bin/sh
# -F: foreground
# -R: create hostkeys as required
# -s: disable password login
# -g: disable password login for root
exec dropbear -FRsg
