# initrd/rootfs/etc/service/ssh/run
#!/bin/sh
#exec > >(logger -t ssh) 2>&1
exec dropbear -F -sg -R
