# ramfs/rootfs/etc/service/ntp/run
#!/bin/sh
exec > >(logger -t ntp) 2>&1
sv start dhcp || exit 1
ntpd -q -p pool.ntp.org
exec ntpd -n
