# initrd/rootfs/etc/service/cron/run
#!/bin/sh
exec > >(logger -t cron) 2>&1
exec crond -fS
