# initrd/rootfs/etc/service/ntp/run
#!/bin/sh
exec > >(logger -t ntp) 2>&1
exec ntpd -n
